# systemd/login-backup.service
# Systemd user service
# Runs the backup script automatically at login
# Used to trigger the Nextcloud incremental backup workflow

[Unit]
Description=Nextcloud Incremental Backup at Login
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sleep 60
ExecStart=%h/login-backup.sh

[Install]
WantedBy=graphical-session.target
